# tools/upnp/scripts/routerDiscovery@.service
[Unit]
Description=Router Discovery Service %i

[Service]
Type=simple
ExecStart=/usr/bin/routerDiscovery %i
