# upnp/upnpdiscover_interface@.service
[Unit]
Description=UPNP Discovery Service %i

[Service]
ExecStart=/usr/bin/upnpdiscover %i
